# nasa-daily-color-avg@pi.timer
[Unit]
Description=Calls the nasa-daily-color-avg service

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=yes

[Install]
WantedBy=timers.target
